# ops/ansible/roles/legislation_explorer/templates/systemd/legislation-explorer.service.j2
[Unit]
Description=OpenFisca Legislation Explorer

[Service]
ExecStart=/usr/bin/npm start
WorkingDirectory={{ legislation_explorer_source_dir | quote }}
User={{ unix_user_name }}
Group={{ unix_group_name }}

[Install]
WantedBy=multi-user.target
